# service/docker/run
#!/bin/sh
/usr/local/bin/docker -d -s vfs 2>&1
